# tools/autonomy/com.clawdbot.autonomy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.clawdbot.autonomy</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/Users/ksato/workspace/projects/openclaw/tools/autonomy/daemon.cjs</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/ksato/workspace/projects/openclaw</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>AUTONOMY_CONFIG</key>
    <string>/Users/ksato/workspace/projects/openclaw/tools/autonomy/config.json</string>
    <key>CODEX_PATH</key>
    <string>/Users/ksato/workspace/shared/_codex</string>
    <key>CLAWDBOT_CLI_TOOLS_DISABLED</key>
    <string>0</string>
    <key>PATH</key>
    <string>/Users/ksato/.npm-global/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/clawdbot-autonomy.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/clawdbot-autonomy.err.log</string>
</dict>
</plist>
